# .config/systemd/user/flrig.service
[Unit]
Description=flrig
PartOf=graphical-session.target

[Service]
ExecStart=flrig

[Install]
WantedBy=xsession.target

