# .config/systemd/user/rigctld.service
[Unit]
Description=rigctld
Wants=flrig.service

[Service]
ExecStart=rigctld -m 4

[Install]
WantedBy=default.target
