# .config/systemd/user/yubikey-touch-detector.service
[Unit]
Description=Yubikey Touch Detector

[Service]
ExecStart=%h/go/bin/yubikey-touch-detector --libnotify

[Install]
WantedBy=default.target

